# contrib/memd-sweep.service
[Unit]
Description=memd inbox sweep

[Service]
Type=oneshot
ExecStart=memd sweep
Nice=10
